# ansible/roles/software/ckan/catalog/ckan-app/files/etc/cron.d/ckan_qa_worker
*/5 * * * * root supervisorctl start celery-mem-check > /dev/null
