# ansible/roles/software/ckan/catalog/ckan-app/files/etc/cron.d/ckan_main_worker
*/3 * * * * root supervisorctl start harvest-run >> /var/log/harvester_run.log 2>&1
