# ansible/roles/software/jekyll/templates/etc/supervisor/conf.d/jekyll.conf.j2
;following entries are generated by jinja2 template
{% for app in jekyll_apps %}

[program:jekyll-{{app.name}}]
command=/usr/local/rvm/gems/ruby-{{ruby_version}}@jekyll/wrappers/jekyll-rebuilder --port {{app.port}} --secret {{ app.secret_key }} --config {{app.config}} --branch {{app.branch|default('master')}} --s3bucket {{app.s3_bucket}}
directory=/root/{{app.name}}
startretries=10000             ; max # of serial start failures (default 3)
stdout_logfile=/var/log/jekyll-{{app.name}}.log
redirect_stderr=true
autostart=true

{% endfor %}
