# ansible/roles/software/catalog/harvest/files/etc/cron.d/ckan-primary
*/3 * * * * root ckan --plugin=ckanext-harvest harvester run >> /var/log/harvester_run.log 2>&1
