# ansible/roles/software/catalog/harvest/files/etc/cron.d/ckan_qa_worker
*/5 * * * * root /usr/bin/celery-mem-check.sh >> /var/log/celery-mem-check.log 2>&1
 
