# ansible/roles/software/ckan/inventory/templates/etc/supervisor/conf.d/gunicorn-web-app.conf
[program:gunicorn-web-app]
command={{ project_source_new_code_path }}/bin/gunicorn --paste /etc/ckan/production.ini -b localhost:5000
directory={{ project_source_new_code_path }}
user=www-data
group=www-data
stdout_logfile={{ inventory_gunicorn_error_log }}
autostart=false
autorestart=false
redirect_stderr=true
