# ansible/roles/software/ckan/inventory/templates/etc/supervisor/conf.d/gunicorn-web-app.conf
[program:gunicorn-web-app]
command={{ project_source_new_code_path }}/bin/gunicorn --paste /etc/ckan/production.ini -b localhost:5000
directory={{ project_source_new_code_path }}
stdout_logfile={{ inventory_gunicorn_error_log }}
autostart=true
autorestart=true
redirect_stderr=true
