# ansible/roles/software/ckan/catalog/harvest/files/etc/cron.d/ckan_qa_worker
*/5 * * * * root supervosorctl start celery-mem-check > /dev/null
