# ansible/roles/software/ckan/catalog/harvest/files/etc/cron.d/ckan_misc_worker
@monthly root supervisorctl start ckan-metrics-csv > /dev/null

#1 1 * * * root supervisorctl start ckan-clean-deleted
1 3 * * * root supervisorctl start ckan-tracking-update > /dev/null
#30 23 * * * root supervisorctl start qa-update-sel:*
#30 7 * * * root supervisorctl stop qa-update-sel:*
#0 6 * * * root supervisorctl start ckan-report
#0 5 * * 0 root supervisorctl start ckan-db-solr-sync
#30 5 * * * root supervisorctl start ckan-harvest-job-cleanup
#30 4 * * * root supervisorctl start ckan-combine-feeds
#2 2 * * * root supervisorctl start ckan-export-csv
0 5 * * * root supervisorctl start ckan-sitemap > /dev/null

## jsonl export on 3rd day of each month at 2:01 AM
#1 2 3 * * root supervisorctl start ckan-jsonl-export
