# ansible/roles/software/catalog/ckan-harvest/files/etc/cron.d/ckan-harvest
*/3 * * * * root ckan --plugin=ckanext-harvest harvester run -c /etc/ckan/production.ini >> /var/log/harvester_run.log 2>&1
