# ansible/roles/software/catalog/ckan-harvest/files/etc/cron.d/supervisor
@reboot root /sbin/initctl start supervisor
