# ansible/roles/software/ckan/solr/files/etc/cron.d/cleanlog
0 * * * * root find /var/log/tomcat6/ -type f -mtime +2 -delete