# ansible/roles/software/jekyll/files/etc/init/supervisor.conf
start on stopped rc RUNLEVEL=[345]
stop on starting shutdown

respawn
exec /usr/local/bin/supervisord
