# ansible/roles/software/datagov-deploy-fgdc2iso/files/etc/cron.d/clean-tmp-files
# delete those tomcat tmp file older than 2 days
0 0 * * * root find /tmp/ -type f -user tomcat -name "tmp*-*-*-*-*" -mtime +2 -delete